# crontab
SHELL=/bin/sh
MAILTO=root

03 5 * * *      git  git --git-dir=/home/services/git/repositories/Refs.git gc --aggressive
*/15 * * * *    git  git --git-dir=/home/services/git/repositories/Refs.git gc
